# launchd/com.cmux.autogroup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.cmux.autogroup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>__HOME__/.local/bin/cmux-autogroup</string>
		<string>watch</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Applications/cmux.app/Contents/Resources/bin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>__HOME__/.local/state/cmux-autogroup/launchd.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/.local/state/cmux-autogroup/launchd.log</string>
</dict>
</plist>
